# etc/NetworkManager/dispatcher.d/90-sbar
#!/bin/sh

case "$2" in
	up|down) kill -RTMIN+2 $(pidof -x sbar);;
esac
